# tools/wiretap-server/can-interface.service
[Unit]
Description=CAN Interface Setup
DefaultDependencies=no
After=local-fs.target
Before=network-pre.target sysinit.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Standard CAN at 500 kbps - adjust bitrate as needed
ExecStart=/sbin/ip link set can0 up type can bitrate 500000

# For CAN FD, use this instead:
# ExecStart=/sbin/ip link set can0 up type can bitrate 500000 dbitrate 2000000 fd on

# For multiple interfaces, add additional ExecStart lines:
# ExecStart=/sbin/ip link set can1 up type can bitrate 500000

ExecStop=/sbin/ip link set can0 down

[Install]
WantedBy=sysinit.target
